# focal/pine64/overlay/lib/systemd/system/ssh.service.d/sshd-config.conf
[Unit]
RequiresMountsFor=/etc/sshd

[Service]
# ssh-generate-hostkeys has to be run before config testing, otherwise testing
# will bail out because it also check for host keys. Because the test is defined
# in the main unit config, we have to clear the list first.
# This is necessary, because the RootFS build removes the host keys by calling
# livebuild's remove-openssh-server-host-keys hook, because we do not want all
# devices to use the same host keys.
ExecStartPre=
ExecStartPre=/usr/bin/ssh-generate-hostkeys
ExecStartPre=/usr/sbin/sshd -t

