# focal/pine64/overlay/lib/systemd/system/lomiri-location-service.service.d/pine.conf
[Service]
ExecStart=
ExecStart=/usr/bin/lomiri-location-serviced --bus system --provider gpsd::Provider --provider mls::Provider