# focal/pine64/overlay/lib/systemd/system/ofono.service.d/pine.conf
[Service]
ExecStartPost=/usr/share/ofono/scripts/enable-modem